# svscode.service
[Unit]
Description=VSCode server
After=network.target

[Service]
Type=simple
ExecStart=/bin/code-server.sh

[Install]
WantedBy=multi-user.target
